# cms/city_forecast.cron
0 */3 * * * cd /home/app/web && python3 manage.py generate_forecast
