# cms/wdqms.cron
0 0 * * * cd /home/app/web && python3 manage.py fetch_transmissions > /proc/1/fd/1 2>/proc/1/fd/2
